# tts-program-server.service.in
[Unit]
Description=Text-to-Speech Program Server
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=@USER@
Group=@GROUP@
Environment=XDG_RUNTIME_DIR=/run/user/@USERID@
ExecStart=@PROGRAM@
WorkingDirectory=@HOME@  
Environment="PATH=@PATH@:@HOME@/.local/bin"

[Install]
WantedBy=multi-user.target

